# kafka/supervisor/kafka-nozk.conf
[program:kafka]
command=/usr/bin/start-kafka-nozk.sh
autostart=true
autorestart=true